# crontab
0 7 * * * django /code/venv/bin/python /code/src/manage.py execute_action --action=EVERY_DAY &> /logs/cron_action_every_day.log
